# RapidApiDataSourcePipeline/crontab
#START CRON JOB 3AM
0 3 * * * python3 /app/RapidAPIFramework.py >>/tmp/out.log 2>/tmp/err.log
#END CRON JOB